# deployment/supervisor/core-server.conf
[program:storeops-core]
command=/opt/storeops/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8088
directory=/opt/storeops/apps/core-server
autostart=true
autorestart=true
stderr_logfile=/var/log/storeops/core.err.log
stdout_logfile=/var/log/storeops/core.out.log
